# crostini/home/.config/systemd/user/cros-garcon.service.d/override.conf
[Service]
Environment="GTK_IM_MODULE=fcitx"
Environment="XMODIFIERS=@im=fcitx"
Environment="QT_IM_MODULE=fcitx"
Environment="QT4_IM_MODULE=fcitx"
