# portable/home/.config/systemd/user/mbsync.target
[Unit]
Description=Mailbox synchronization target
Wants=mbsync.service
After=mbsync.service
