# portable/home/.config/systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=2m
OnUnitInactiveSec=4m
Unit=mbsync.service

[Install]
WantedBy=timers.target
